# etc/supervisor/conf.d/foxserver_complete.conf
[program:foxserver_complete]
command=/root/bin/check_foxserver.sh
redirect_stderr=true
autorestart=false
stdout_logfile = /var/log/supervisor/%(program_name)s.log
